# systemd/update-cloudflare-ufw.service
[Unit]
Description=Update UFW rules with latest Cloudflare IPs
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-cloudflare-ufw.sh
